# config_files/backups/20251125_091350/yasrl-api.service
[Unit]
Description=YASRL RAG API Server
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=notify
User=deploy
Group=www-data
RuntimeDirectory=yasrl-api
RuntimeDirectoryMode=755
WorkingDirectory=/srv/library2
Environment=PATH=/srv/library2/venv/bin
Environment=PYTHONPATH=/srv/library2
EnvironmentFile=/srv/library2/.env
ExecStart=/srv/library2/venv/bin/gunicorn --config /srv/library2/gunicorn.conf.py src.yasrl.api:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=on-failure
RestartSec=10

# Security settings
NoNewPrivileges=yes
ProtectSystem=yes
ProtectHome=yes
ReadWritePaths=/srv/library2 /run/yasrl-api /tmp /var/log/yasrl
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
